# .config/systemd/user/mount-sshfs@Videos.service.d/mount.conf
[Unit]
Requires=mount-veracrypt@ssh.service
After=mount-veracrypt@ssh.service

[Service]
Environment="TARGET=portnoy4prez:/mnt/data/home/Videos"
Environment="MOUNTPOINT=/mnt/data/Videos"
