# .config/systemd/user/autogit.service
[Unit]
Description=Sync git repos

[Service]
Type=oneshot
ExecStart=%h/.bin/autogit "$REPOS"

[Install]
WantedBy=multi-user.target
