# .config/systemd/user/autogit.timer
[Unit]
Description=Run auto every hour

[Timer]
Persistent=true
OnCalendar=hourly

[Install]
WantedBy=timers.target