# .config/systemd/user/batmon.service
[Unit]
Description=Batmon service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/batmon

[Install]
WantedBy=multi-user.target