# .config/systemd/user/batmon.timer
[Unit]
Description=Run batmon every minute

[Timer]
OnBootSec=0min
Persistent=true
OnCalendar=*:0/1

[Install]
WantedBy=timers.target