# .config/systemd/user/conky_pacman.timer
[Unit]
Description=Conky Pacman Stats Update Timer

[Timer]
OnBootSec=0s
OnUnitActiveSec=1m

[Install]
WantedBy=timers.target