# .config/systemd/user/mbsync.timer
[Unit]
Description=Mailbox synchronization timer

[Timer]
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target