# .config/systemd/user/mount-veracrypt@ssh.service.d/mount.conf
[Service]
Environment="BW_NAME=Veracrypt (ssh)"
Environment="VOLUME=/mnt/data/Documents/crypt/ssh-config"
Environment="MOUNTPOINT=%h/.ssh"