# .config/systemd/user/mount-veracrypt@gnupg.service.d/mount.conf
[Service]
Environment="BW_NAME=Veracrypt (gpg)"
Environment="VOLUME=/mnt/data/Documents/crypt/gnupg"
Environment="MOUNTPOINT=$GNUPGHOME"