# .config/systemd/user/autogit.service.d/autogit.conf
[Service]
Environment="REPOS=\"~/Org ~/Org/journal\""
